# systemd/findmejobs-review-import.service
[Unit]
Description=FindMeJobs OpenClaw review import

[Service]
Type=oneshot
User=findmejobs
Group=findmejobs
WorkingDirectory=/opt/findmejobs
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/findmejobs/.venv/bin/findmejobs review import-results --app-config-path /etc/findmejobs/app.toml --profile-path /etc/findmejobs/profile.yaml --sources-path /etc/findmejobs/sources.yaml

[Install]
WantedBy=multi-user.target
