# systemd/findmejobs-review-import.timer
[Unit]
Description=FindMeJobs OpenClaw review import schedule

[Timer]
OnCalendar=*:10/10
Persistent=true

[Install]
WantedBy=timers.target
