# systemd/findmejobs-review-export.timer
[Unit]
Description=FindMeJobs review packet export schedule

[Timer]
OnCalendar=*:4/15
Persistent=true

[Install]
WantedBy=timers.target
