# systemd/findmejobs-digest.timer
[Unit]
Description=FindMeJobs digest delivery schedule

[Timer]
OnCalendar=*-*-* 08:30:00
Persistent=true

[Install]
WantedBy=timers.target
