# systemd/findmejobs-review-export.service
[Unit]
Description=FindMeJobs review packet export
After=findmejobs-rank.service

[Service]
Type=oneshot
User=findmejobs
Group=findmejobs
WorkingDirectory=/opt/findmejobs
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/findmejobs/.venv/bin/findmejobs review export --app-config-path /etc/findmejobs/app.toml --profile-path /etc/findmejobs/profile.toml --sources-dir /etc/findmejobs/sources.d

[Install]
WantedBy=multi-user.target
