# systemd/findmejobs-digest.service
[Unit]
Description=FindMeJobs digest delivery
After=findmejobs-review-import.service

[Service]
Type=oneshot
User=findmejobs
Group=findmejobs
WorkingDirectory=/opt/findmejobs
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/findmejobs/.venv/bin/findmejobs digest send --app-config-path /etc/findmejobs/app.toml --profile-path /etc/findmejobs/profile.toml --sources-dir /etc/findmejobs/sources.d

[Install]
WantedBy=multi-user.target
